# crontab
3 * * * * /usr/bin/fetchmail -a --verbose


